# deployment/dbmanager-api.service
[Unit]
Description=DBManager REST API Server
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/dbmanager
ExecStart=/usr/bin/python3 /path/to/dbmanager/api_server.py
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/dbmanager-api.log
StandardError=append:/var/log/dbmanager-api-error.log

[Install]
WantedBy=multi-user.target
